# distribution/debian/etc/service/oreka-api/run
#!/bin/bash
exec 2>&1

/usr/local/bin/oreka-api